# deploy/systemd/siren-driver.service
[Unit]
Description=Siren Driver webhook service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{SERVICE_USER}}
Group={{SERVICE_USER}}
SupplementaryGroups=gpio
WorkingDirectory={{HOME_DIR}}
Environment=HOME={{HOME_DIR}}
Environment=PYTHONPATH={{APP_DIR}}/src
EnvironmentFile={{ENV_FILE}}
ExecStart={{VENV_DIR}}/bin/uvicorn siren_driver.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5
UMask=0027
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectClock=yes
ProtectHostname=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
ReadWritePaths={{HOME_DIR}} {{CONFIG_DIR}}

[Install]
WantedBy=multi-user.target
